# supervisord.base.conf
[supervisord]
inodaemon=false

[program:php]
command=php-fpm

[program:lottery-worker]
command=/usr/sbin/nginx
